# src/clean-pacman-cache.timer
[Unit]
Description=Clear the pacman cache every week

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
